# prediction-service/crontab
# Run training every day at 3:00 AM
0 3 * * * python /app/scripts/train.py >> /var/log/cron.log 2>&1
# Important: Leave a blank line at the end of the file
